# blob_server.conf
start on startup
exec sh /home/start.sh

#place this file in /etc/init/