# simple_server.conf
start on startup
exec node /home/ftpuser/index.js

#place this file in /etc/init/