# simple_server.conf
start on startup
exec node /home/ftpuser/simple_server.js

#place this file in /etc/init/